# deployment/supervisord/conf.d/welt2000.conf
[program:welt2000]
user=welt2000
command=/home/welt2000/venv/bin/uwsgi -s /tmp/uwsgi-welt2000.sock --chmod-socket=664 --chown-socket=welt2000:www-data --module welt2000 --callable app
directory=/home/welt2000/src
autostart=true
autorestart=true
stdout_logfile=/home/welt2000/log/uwsgi.log
redirect_stderr=true
stopsignal=QUIT
virtualenv=
